# autostart/.config/autostart/dunst.desktop
[Desktop Entry]
Type=Application
Name=dunst
Icon=dunst
Comment=Notification daemon
Path=/usr/bin
Exec=/usr/bin/dunst
Terminal=false
